# deploy/self-hosting/systemd/lem-signaling.service
[Unit]
Description=Lem Signaling Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=lem
Group=lem
WorkingDirectory=/opt/lem/cloud/signaling

# Environment
Environment="PATH=/opt/lem/.local/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/opt/lem/cloud/signaling/.env

# Start command using uvicorn with uv
ExecStart=/opt/lem/.local/bin/uv run uvicorn app.main:app \
    --host ${HOST} \
    --port ${PORT} \
    --workers 4 \
    --loop uvloop \
    --log-level info

# Restart policy
Restart=always
RestartSec=5
StartLimitInterval=0

# Resource limits
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/lem/cloud/signaling

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=lem-signaling

[Install]
WantedBy=multi-user.target
